# axon.service
[Unit]
Description=Axon API Gateway
Wants=syslog.target network-online.target
After=network.target network-online.target

[Service]
Type=simple
User=akagi201
Group=akagi201
ExecStart=/usr/local/bin/axon --config /etc/axon/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/axon
ReadOnlyPaths=/etc/axon

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Environment variables
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Working directory
WorkingDirectory=/var/lib/axon

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
